# roles/mss_nfqws/files/mss-nfqws.service.d/netns.conf
[Unit]
JoinsNamespaceOf=lan_guest_ns.service
BindsTo=lan_guest_ns.service
Requires=lan_guest_ns.service
After=lan_guest_ns.service

[Service]
PrivateNetwork=Yes
PrivateMounts=No
