# roles/bind/files/named.service.d/breakdown.conf
[Service]
ExecStartPre=/opt/mss/share/multiwg-handler named
